# ai-stack/systemd/kiwix-update.timer
[Unit]
Description=Monthly Kiwix ZIM update check
Requires=kiwix-update.service

[Timer]
OnCalendar=monthly
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
